# root.cron
0,30 * * * *	. /etc/docker-gc.env && /docker-gc
